# systemd/user/libation-s3.service
[Unit]
Description=Mount S3 (rclone) for Libation
After=network-online.target
Wants=network-online.target
PartOf=libation.service
BindsTo=libation.service

[Service]
Type=simple
Slice=app-libation.slice
EnvironmentFile=%h/.config/libation-systemd/env

# Let the user know we’re starting the mount
ExecStartPre=-/usr/bin/notify-send --app-name=Libation --icon=libation "Mounting storage…" "Connecting to S3 and preparing mount"

# Best-effort unmount if something is already there
ExecStartPre=-/bin/fusermount -u -z ${MOUNTPOINT}
# Make sure the mountpoint exists
ExecStartPre=/usr/bin/mkdir -p ${MOUNTPOINT}

# Run rclone in the foreground
ExecStart=/usr/bin/rclone mount :s3:${BUCKET} ${MOUNTPOINT} \
  --s3-no-check-bucket \
  --uid %U --gid %G --umask 002 \
  --vfs-cache-mode writes \
  --vfs-cache-poll-interval 30m \
  --poll-interval 0 \
  --dir-cache-time 24h \
  --attr-timeout 24h \
  --use-server-modtime \
  --bwlimit ${UP_LIMIT}:${DOWN_LIMIT} \
  --cache-dir /tmp/libation \
  --exclude "/.Trash-*/**" \
  --rc

# If MIGRATION_NEEDED=true
ExecStartPost=/bin/sh -c '/usr/share/libation-systemd/migrate-s3.sh 2>&1 &'

# Clean unmount on stop
ExecStop=/bin/fusermount -u -z ${MOUNTPOINT}

Restart=on-failure
RestartSec=5s
RemainAfterExit=no

TimeoutStartSec=30s
TimeoutStopSec=10s
